# scripts/scripts/newsboat/com.newsboat.refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.newsboat.refresh</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/newsboat</string>
        <string>--refresh-on-start</string>
        <string>--execute</string>
        <string>quit</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/newsboat-refresh.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/newsboat-refresh.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
